# deploy/twbshop-automations.service
[Unit]
Description=TWBshop automations auto-dispatch runner (sends firing recipes to configured targets)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/TWBshop
Environment=TWBSHOP_ENV=prod
# ExecStart python path matches the other twbshop-* units (verified against twbshop-wizard.service).
ExecStart=/root/venv/bin/python3 run_automations.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
