# deploy/twbshop-retail.service
[Unit]
Description=TWBshop Retail Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/TWBshop
ExecStart=/root/venv/bin/python3 run_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
