# deploy/twbshop-retail.service
[Unit]
Description=TWBshop Retail Bot
After=network.target

[Service]
Type=simple
User=twbshop
WorkingDirectory=/home/twbshop/TWBshop
ExecStart=/usr/bin/python3 run_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
